# curator/supervisord.conf
[supervisord]
nodaemon=true

[program:t2c2curator]
command=/bin/bash -c "/clowder/curator-start.sh"
